# config/systemd/user/calcurse-caldav.service
[Unit]
Description=calcurse-caldav job

[Service]
ExecStart=/home/fran/dotfiles/scripts/calcurse-caldav-fran
