# config/systemd/user/poll-mail.timer
[Unit]
Description=Poll mail every 15 minutes

[Timer]
OnBootSec=15s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
