# private_dot_config/private_autostart/bitwarden.desktop
[Desktop Entry]
Type=Application
Version=2024.8.1
Name=Bitwarden
Comment=Bitwarden startup script
Exec=bitwarden-desktop
StartupNotify=false
Terminal=false